# rootfs/etc/s6-overlay/s6-rc.d/worker/run
#!/command/with-contenv sh
set -eu

cd /rails

: "${DB_HOST:=127.0.0.1}"
: "${DB_PORT:=5432}"
: "${POSTGRES_USER:=sure_user}"
: "${POSTGRES_PASSWORD:=internal_sure_pass}"
: "${POSTGRES_DB:=sure_production}"
: "${REDIS_URL:=redis://127.0.0.1:6379/1}"

export DB_HOST DB_PORT POSTGRES_USER POSTGRES_PASSWORD POSTGRES_DB REDIS_URL

echo "Waiting for Sure web service to finish bootstrapping..."
until curl -fsS http://127.0.0.1:3000/up >/dev/null 2>&1; do
  sleep 2
done

# Launch Background Worker
exec bundle exec sidekiq
