# rootfs/etc/services.d/sure-worker/run
#!/command/execlineb -P
exec bundle exec sidekiq
